# systemd/sshflingd.service
[Unit]
Description=SSHFling temporary SSH certificate issuer
Documentation=file:/usr/share/doc/sshfling/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sshfling/sshflingd.env
ExecStart=/usr/bin/sshfling serve --listen ${SSHFLING_LISTEN} --ca-key ${SSHFLING_CA_KEY} --allowed-principal ${SSHFLING_ALLOWED_PRINCIPALS} --default-seconds ${SSHFLING_DEFAULT_SECONDS} --max-seconds ${SSHFLING_MAX_SECONDS} --session-wrapper ${SSHFLING_SESSION_WRAPPER}
User=sshflingd
Group=sshflingd
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/sshfling

[Install]
WantedBy=multi-user.target
